# webapp.service
[Unit]
Description= Webapp

[Service]
ExecStart=/usr/bin/npm run dev
Restart=always
User=ec2-user
WorkingDirectory=/home/ec2-user/webapp

[Install]
WantedBy=multi-user.target
